# deploy/supervisor.d/elastic-searchd.ini
[program:es] 
command=cd /opt/elasticsearch-8.12.0/bin && ./elasticsearch
stdout_logfile=/data/es/logs/elastic-searchd.log
autostart=true
autorestart=true
startsecs=5
user=es-user
redirect_stderr=true
stopasgroup=true
killasgroup=true
