# deploy/supervisor.d/minio.ini
[program:minio] 
command=/etc/supervisord.d/minio.sh
stdout_logfile=/data/logs/minio.log
autostart=true
autorestart=true
startsecs=5
redirect_stderr=true
stopasgroup=true
killasgroup=true
