# deploy/supervisor.d/rag-broker.ini
[program:rag-broker]
command=/etc/supervisord.d/broker-srv.sh
stdout_logfile=/data/logs/rag-broker-srv.log
autostart=true
autorestart=true
startsecs=5
redirect_stderr=true
stopasgroup=true
killasgroup=true
