# deploy/supervisor.d/rag-executor.ini
[program:rag-executor]
command=/etc/supervisord.d/executor-srv.sh
stdout_logfile=/data/logs/rag-executor-srv.log
autostart=true
autorestart=true
startsecs=5
redirect_stderr=true
stopasgroup=true
killasgroup=true
