# deploy/supervisor.d/rag-api.ini
[program:rag-api]
command=/etc/supervisord.d/api-srv.sh
stdout_logfile=/data/logs/rag-api-srv.log
autostart=true
autorestart=true
startsecs=5
redirect_stderr=true
stopasgroup=true
killasgroup=true
